# spring-cloud-eureka/eureka.supervisor.ini
[supervisord]
nodaemon=true

[supervisorctl]
port=127.0.0.1:9001

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:eureka]
directory=/Users/sh/repos/yrayjava/spring-cloud-eureka
command=bash -c "mvn -Dcustom.mode=dev clean spring-boot:run"
autostart=true
autorestart=true
stopsignal=QUIT
stopwaitsecs=2   ;执行SIGKILL之前对多等待时间(s)，目前本项目没有什么重要的事务需要必须执行完，所以设置的2s
stopasgroup=true
killasgroup=true
stopsignal=QUIT
process_name=%(program_name)s
stdout_logfile=/tmp/logs/tomcat.log
stderr_logfile=/tmp/logs/tomcat.err